# ansible/ansible-apps/files/api.service
[Unit]
Description=API Service
After=network.target

[Service]
ExecStart=/tmp/api
Restart=always
User=root

[Install]
WantedBy=multi-user.target
